# supervisor-config.conf
[program:transcoder]
command=/usr/local/bin/transcoder.py
autorestart=true
